# cron/lottery.cron
* * * * * /usr/local/nginx/lottery/cron/update_lottery.php
